# statigo.service.example
[Unit]
Description=Statigo Web Application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=statigo
Group=statigo
WorkingDirectory=/opt/statigo
EnvironmentFile=/opt/statigo/.env
ExecStart=/opt/statigo/statigo
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=statigo

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/statigo/data/cache /opt/statigo/data

# Resource limits
LimitNOFILE=65536
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
